# source/rootfs/scripts/rootfs/var/lib/dpkg/info/rsyslog.conffiles
/etc/apparmor.d/usr.sbin.rsyslogd
/etc/init.d/rsyslog
/etc/logcheck/ignore.d.server/rsyslog
/etc/logrotate.d/rsyslog
/etc/rsyslog.conf
